# packaging/calendar-serviced.service
[Unit]
Description=Calendar Service daemon
After=contacts-service.service
Requires=contacts-service.service

[Service]
ExecStart=/usr/bin/calendar-serviced
Type=dbus
BusName=ALARM.acalendar-service

[Install]
WantedBy=default.target
