# packaging/calendar-serviced.service
[Unit]
Description=Calendar Service daemon
After=contacts-service.service
Requires=contacts-service.service

[Service]
Type=simple
ExecStart=/usr/bin/calendar-serviced

[Install]
WantedBy=default.target
